# systemd/void-herald.service
# ═══════════════════════════════════════════════════════════════
# void-herald.service  —  copy to /etc/systemd/system/void-herald.service
# ═══════════════════════════════════════════════════════════════
# Before installing, replace:
#   YOUR_USERNAME  → void-herald  (or whatever system user you created)
#   YOUR_HOME      → the directory that holds Claude Code's config
#   BOT_PATH       → absolute path to the void-herald checkout
#
# Create the service account first:
#   sudo useradd --system --no-create-home --shell /usr/sbin/nologin void-herald
#
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now void-herald

[Unit]
Description=void-herald Discord Bot (Node.js)
After=network.target ngrok.service

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=BOT_PATH
ExecStart=/usr/bin/node index.js
Restart=on-failure
RestartSec=5
EnvironmentFile=BOT_PATH/.env
Environment=HOME=YOUR_HOME

StandardOutput=journal
StandardError=journal

# ── Hardening ────────────────────────────────────────────────
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
CapabilityBoundingSet=

[Install]
WantedBy=multi-user.target
